# ansible/roles/cri/templates/crio.service
[Service]
Delegate=yes
Slice=podruntime.slice
